# switchbot-data-collector.service
[Unit]
Description=SwitchBot Power Data Collection
After=switchbot-power-monitor.service
Requires=switchbot-power-monitor.service

[Service]
Type=oneshot
User=user
WorkingDirectory=/home/user/switchbot-power-monitor
ExecStart=/home/user/switchbot-power-monitor/.venv/bin/python /home/user/switchbot-power-monitor/collect_data.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target